# docker/dev/php/files/etc/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
silent=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:php-fpm]
command=/usr/sbin/php-fpm7.4 -c /etc/php/7.4/fpm/php-fpm.conf -F
stdout_logfile=/tmp/supervisord.log
stdout_syslog=true
stderr_logfile=/tmp/supervisord.log
stderr_syslog=true

